# waybar/hypr/scripts/pacman-db-update.service
[Unit]
Description=Update pacman package database

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Sy --quiet
